# supervisor_bjdns_example.conf
[program:bjdns]
directory = /home/user/bjdns/
command = python3 bjdns3
autostart = true
startsecs = 5
user = user
redirect_stderr = true
stdout_logfile = /home/user/bjdns/bjdns.log
loglevel = debug
